# .config/systemd/user/epbg.service
[Unit]
Description=Download and set desktop background to a random /r/earthporn image

[Service]
Type=oneshot
ExecStart=/bin/bash $HOME/bin/epbg.sh
